# dashboard.service
[Unit]
Description=Send It Trading Dashboard
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/path/to/send-it-trading
ExecStart=/usr/bin/python3 dashboard_api.py
Restart=always
RestartSec=10
EnvironmentFile=/path/to/send-it-trading/.env

[Install]
WantedBy=multi-user.target
