# services/bot.service
# Template — run install.sh to generate configured versions with correct paths
# Or replace REPO_DIR with your absolute repo path and USERNAME with your user

[Unit]
Description=Send It Trading — Intelligence Bot (news/insider/polymarket scanners)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USERNAME
WorkingDirectory=REPO_DIR
ExecStart=REPO_DIR/venv/bin/python REPO_DIR/bot/main.py
Restart=on-failure
RestartSec=30
EnvironmentFile=REPO_DIR/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=send-it-bot

[Install]
WantedBy=multi-user.target
