# services/dashboard.service
# Template — run install.sh to generate configured versions with correct paths
# Or replace REPO_DIR with your absolute repo path and USERNAME with your user

[Unit]
Description=Send It Trading — Live Dashboard
After=network.target

[Service]
Type=simple
User=USERNAME
WorkingDirectory=REPO_DIR
ExecStart=REPO_DIR/venv/bin/python REPO_DIR/dashboard/api.py
Restart=always
RestartSec=10
EnvironmentFile=REPO_DIR/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=send-it-dashboard

[Install]
WantedBy=multi-user.target
