# services/dashboard.service
[Unit]
Description=Send It Trading — Live Dashboard
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/send-it-trading
ExecStart=/usr/bin/python3 %h/send-it-trading/dashboard/api.py
Restart=always
RestartSec=10
EnvironmentFile=%h/send-it-trading/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=send-it-dashboard

[Install]
WantedBy=multi-user.target
