# services/engine.service
# Template — run install.sh to generate configured versions with correct paths
# Or replace REPO_DIR with your absolute repo path and USERNAME with your user

[Unit]
Description=Send It Trading — Options Engine (alpha scoring + execution)
After=network.target send-it-bot.service

[Service]
Type=simple
User=USERNAME
WorkingDirectory=REPO_DIR
ExecStartPre=/bin/sleep 15
ExecStart=REPO_DIR/venv/bin/python REPO_DIR/engine/main_wrapper_simple.py
Restart=always
RestartSec=10
EnvironmentFile=REPO_DIR/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=send-it-engine
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
